# scripts/crontab/test.cron
30 18 12 * * /home/commaai-03/Documents/autosave/20181108.sh
